# deployment/supervisor/myblog-tasks.conf
[program:myblog-tasks]
command=/home/jim/myblog/venv/bin/rq worker myblog-tasks
numprocs=1
directory=/home/jim/myblog
user=jim
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

