# deployment/supervisor/myblog.conf
[program:myblog]
command=/home/jim/myblog/venv/bin/gunicorn -b localhost:8000 -w 4 myblog:app
directory=/home/jim/myblog
user=jim
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

